# .config/cron/crontab
LANG="en_US.UTF-8"
PATH="~/.nix-profile/bin:/usr/bin:/bin"

* * * * * source ~/.profile && bash ~/.config/notmuch/move.bash >> ~/cron_notmuch.log 2>> ~/cron_notmuch_error.log
