# .config/cron/crontab
LANG="en_US.UTF-8"
PATH="~/.nix-profile/bin:/usr/bin:/bin"

*/5 * * * * source ~/.profile && bash ~/.config/notmuch/move.bash >> ~/.cache/cron_notmuch.log 2>> ~/.cache/cron_notmuch_error.log
